# deploy/ntfy-bot.service
[Unit]
Description=Amul Protein Stock Notifier Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/ntfy-bot
ExecStart=%h/ntfy-bot/ntfy_bot
EnvironmentFile=%h/ntfy-bot/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
