# deploy/ntfy-bot.service
[Unit]
Description=Amul Protein Stock Notifier Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ntfy-bot
ExecStart=/opt/ntfy-bot/ntfy_bot
EnvironmentFile=/opt/ntfy-bot/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
